# files/21.02/etc/rc.local
# Copyright (C) 2010 Stephan Enderlein <stephan@freifunk-dresden.de>
# GNU General Public License Version 3

# DON'T put anything here, because it is saved and restored along
# saving configurations
# freifunk must not call scripts from here

exit 0
