# files/14.07/etc/hotplug.d/iface/20-firewall
#firewall reloading disabled; it would delete all rules set by backbone/vpn.
#all interfaces are preset to "up" state, so all rules are set during boot.

