# feeds/15.05/feeds-own/wshaper/files/wshaper.init
#!/bin/sh /etc/rc.common

START=45
start() {
	if [ "$(uci get ddmesh.network.speed_enabled)" = "1" ]; then
		wshaper.htb
	fi
}

stop() {
	wshaper.htb stop
}
