# files/21.02/etc/hotplug.d/ieee80211/01-eventlog
#!/bin/sh

TAG="ddmesh-hotplug-ieee80211"
logger -t "$TAG" "event: net:[$INTERFACE] action:$ACTION type:$DEVTYPE devicename:[$DEVICENAME] devpath:$DEVPATH"

