# files/common/etc/hotplug.d/iface/01-eventlog
#!/bin/sh
logger -t "ddmesh-hotplug" "iface-event: net:[$INTERFACE] action:$ACTION type:$DEVTYPE devicename:[$DEVICENAME] devpath:$DEVPATH"
