# files/common/etc/hotplug.d/iface/01-eventlog
#!/bin/sh
logger -t ddmesh "network-hotplug iface: [$INTERFACE] action:$ACTION type:$DEVTYPE devicename:[$DEVICENAME] devpath:$DEVPATH"
