# files/common/etc/hotplug.d/iface/01-eventlog
#!/bin/sh
logger -t ddmesh "network-hotplug iface: [$INTERFACE] $ACTION type:$DEVTYPE interface:[$DEVICENAME] devpath:$DEVPATH"

