# files/common/etc/hotplug.d/iface/01-eventlog
#!/bin/sh
logger -t ddmesh "network-hotplug iface: network:[$INTERFACE] interface:$DEVICE $ACTION"
