# files/common/etc/hotplug.d/iface/25-dnsmasq
#disable dnsmasq reload when interface setting change
#it would disturb current dhcp leases when tbb+ interfaces or bat+
#change
