# files/common/etc/hotplug.d/iface/51-disable-bmxd-gateway

test -f /tmp/freifunk-running || exit 0

#echo "$ACTION:$INTERFACE:$DEVICE" >>/tmp/x
#deaktivate gateway immediately, when cable is unplugged
[ ifdown = "$ACTION" ] && {
	[ wan = "$INTERFACE" -o lan = "$INTERFACE" ] && /usr/lib/ddmesh/ddmesh-bmxd.sh no_gateway
}

