# files/common/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

TAG="ddmesh-hotplug"
logger -t "$TAG" "ieee80211-event: net:[$INTERFACE] action:$ACTION type:$DEVTYPE devicename:[$DEVICENAME] devpath:$DEVPATH"

