# files/common/etc/hotplug.d/ieee80211/10-wifi-detect
# freifunk
# disable hotplug, it runs asynchroniously to freifunk wifi setup
