# files/common/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

[ "${ACTION}" = "add" ] && {
	logger -s -t "hotplug:" "ieee80211 ADD"
	/usr/lib/ddmesh/ddmesh-hotplug-wifi.sh
}
 
