# files/21.02/etc/hotplug.d/iface/01-eventlog
#!/bin/sh
TAG="ddmesh-hotplug-iface"
logger -t "$TAG" "event: net:[$INTERFACE] action:$ACTION type:$DEVTYPE devicename:[$DEVICENAME] devpath:$DEVPATH"
