# files/common/etc/hotplug.d/net/01-eventlog
#!/bin/sh
TAG="ddmesh-hotplug"

logger -t "$TAG" "net-event: net:[$INTERFACE] $ACTION type:$DEVTYPE interface:[$DEVICENAME] devpath:$DEVPATH"

