# files/common/etc/hotplug.d/net/01-eventlog
#!/bin/sh
TAG="ddmesh-hotplug-net"

logger -t "$TAG" "event: net:[$INTERFACE] action:$ACTION type:$DEVTYPE interface:[$DEVICENAME] devpath:$DEVPATH"

