# files/common/etc/hotplug.d/net/01-eventlog
#!/bin/sh
logger -t ddmesh "network-hotplug net: [$INTERFACE] $ACTION type:$DEVTYPE interface:[$DEVICENAME] devpath:$DEVPATH"

